# php-worker/supervisord.d/lavarel-swoole-9501.conf
[program:laravel-swoole-9501]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/sns/ artisan swoole:start --port=9501
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true